# supervisord.conf
[supervisord]
nodaemon=true

[program:iota]
command=waitress-serve --call 'iota:create_app'
user=iota
